# deploy/routeros.service
[Unit]
Description=RouterOS WebSocket Service - Gerenciamento RouterOS via WebSocket
After=network.target vpnserverio.service
Requires=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/automais.io/vpnserver.io
ExecStart=/root/automais.io/vpnserver.io/venv/bin/python3 /root/automais.io/vpnserver.io/routeros_websocket.py
Restart=always
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=routeros-websocket

# Timeouts
TimeoutStartSec=60
TimeoutStopSec=30
TimeoutAbortSec=10

# Variáveis de ambiente
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/root/automais.io/vpnserver.env

[Install]
WantedBy=multi-user.target

